# ops/com.jokull.icelandic-data-dms.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.jokull.icelandic-data-dms</string>
	<key>Comment</key>
	<string>Dead-man's-switch: alerts if github.com/jokull/icelandic-data source-health stops recording observations. Watches the monitor, not the sources.</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/jokull/clawd/bin/icelandic-data-dms.sh</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<!-- Every 6h. The check tolerates 36h of silence, so this detects a stall
	     within a quarter of that window while costing nothing on healthy runs. -->
	<key>StartInterval</key>
	<integer>21600</integer>
	<key>StandardOutPath</key>
	<string>/Users/jokull/clawd/logs/icelandic-data-dms.launchd.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/jokull/clawd/logs/icelandic-data-dms.launchd.log</string>
</dict>
</plist>
